# profiles/standard/netbook/etc/skel/.config/autostart/cinnamon-screensaver.desktop
[Desktop Entry]
Type=Application
Terminal=false
StartupNotify=true
Exec=cinnamon-screensaver
